# presto-hdp-hive/files/etc/supervisord.d/presto.conf
[program:presto]
command=presto-start.sh
startsecs=2
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/presto/presto-server.log
autostart=true